# supervisord.conf
[supervisord]
nodaemon=true

[program:php8.1-fpm]
command=/usr/local/sbin/php-fpm

[program:nginx]
command=/usr/sbin/nginx
